# infra/gateway/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/gateway/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tailscaled]
command=/usr/local/bin/tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/var/run/tailscale/tailscaled.sock
autostart=true
autorestart=true
stdout_logfile=/var/log/gateway/tailscaled.log
stderr_logfile=/var/log/gateway/tailscaled.err
priority=10

[program:tailscale-up]
command=/usr/local/lib/gateway/tailscale-up.sh
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/var/log/gateway/tailscale-up.log
stderr_logfile=/var/log/gateway/tailscale-up.err
priority=20

[program:caddy]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
stdout_logfile=/var/log/gateway/caddy.log
stderr_logfile=/var/log/gateway/caddy.err
priority=30

[program:squid]
command=/usr/sbin/squid -NYC -f /etc/squid/squid.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/gateway/squid.log
stderr_logfile=/var/log/gateway/squid.err
priority=30

[program:cert-renew]
command=/usr/local/lib/gateway/cert-renew.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/gateway/cert-renew.log
stderr_logfile=/var/log/gateway/cert-renew.err
priority=40
